# launchd/com.stone.intraday.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stone.intraday</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/lindeng/Stone/.venv/bin/python</string>
        <string>/Users/lindeng/Stone/scripts/daily_analysis.py</string>
        <string>--mode</string>
        <string>intraday</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>14</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/lindeng/Stone</string>
    <key>StandardOutPath</key>
    <string>/Users/lindeng/Stone/logs/intraday.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/lindeng/Stone/logs/intraday.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/lindeng/Stone/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
